# systemd/.config/systemd/user/polybar-updates.timer
[Unit]
Description=Check for system updates for polybar

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target