# systemd/.config/systemd/user/waybar-updates.service
[Unit]
Description=Fetch package updates for Waybar

[Service]
Type=oneshot
ExecStart=%h/.config/waybar/update-checker.sh